# scripts/systemd/sss2-frontend.service
# SSS2-GUI Svelte frontend (ui) served by the adapter-node Node server.
# Install:
#   sudo cp scripts/systemd/sss2-frontend.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now sss2-frontend.service
# Edit User=, paths, and PORT to match your Pi. `node build` runs the output of `npm run build`.

[Unit]
Description=SSS2-GUI Svelte frontend (ui, adapter-node)
After=network-online.target sss2-backend.service
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/sss2-gui-v2/ui
Environment=HOST=0.0.0.0
Environment=PORT=3000
# If a reverse proxy terminates TLS for the public domain, set ORIGIN so SvelteKit
# builds correct absolute URLs, e.g.:
# Environment=ORIGIN=https://sss2.example.com
ExecStart=/usr/bin/node build
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
